# crowdsec-dyanmic-allowlist.crontab

*/10  * * * * chdir /opt/security-dynamic-allowlist  && ./reload-cs.ts
