# crowdsec-dynamic-allowlist.crontab

*/10  * * * * chdir /opt/crowdsec-let-meme-in  && ./reload-cs.ts
